# roles/preflight/templates/90-wifi-usb-rename-adapter.rules.j2
ACTION=="add", SUBSYSTEM=="net", ATTR{address}=="00:c0:ca:b0:7d:dd", RUN+="{{ logged_user_home }}/.local/bin/rename-wifi.sh"
